# supervisord.conf
[supervisord]
nodaemon=true                ; start in foreground 

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"

[program:mysqld]
#command = /usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe
command=/start-mysqld.sh
numprocs=1
autostart=true
autorestart=true
